# templates/raspberrypi_exporter.service.j2
[Unit]
Description=Prometheus exporter for Raspberry Pi metrics
Wants=raspberrypi_exporter.timer

[Service]
Type=oneshot
ExecStart={{ raspberrypi_exporter_script_path }}/raspberrypi_exporter
